# launchd/com.trade-analysis.skill-improvement.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.trade-analysis.skill-improvement</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/takueisaotome/PycharmProjects/claude-trading-skills/scripts/run_skill_improvement.sh</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/takueisaotome/PycharmProjects/claude-trading-skills</string>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>5</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/Users/takueisaotome/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>/Users/takueisaotome</string>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/takueisaotome/PycharmProjects/claude-trading-skills/logs/launchd_skill_improvement.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/takueisaotome/PycharmProjects/claude-trading-skills/logs/launchd_skill_improvement_error.log</string>

	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
